# profile.d/50homebrew.sh
if which brew > /dev/null; then
    source_directory homebrew
fi
